# code-deploy/personal-site-netcore.conf
# This is an upstart (http://upstart.ubuntu.com/) script
# to run the node.js server on system boot and make it
# manageable with commands such as
# 'start node-app' and 'stop node-app'
#
# This script is to be placed in /etc/init to work with upstart.
#
# Internally the 'initctl' command is used to manage:
# initctl help
# initctl status node-app
# initctl reload node-app
# initctl start node-app

# Got this mostly from http://kvz.io/blog/2009/12/15/run-nodejs-as-a-service-on-ubuntu-karmic/
# which was referenced by http://stackoverflow.com/questions/11084279/node-js-setup-for-easy-deployment-and-updating

description "Launch my personal website in go and keep it running indefinitely."
author      "Ethan Young"
version "1.0"

# used to be: start on startup
# until we found some mounts weren't ready yet while booting:

start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 10 5

script
    # Log output while starting the site
    exec >/var/log/personal-site-go-start.log 2>&1

    # Not sure why $HOME is needed, but we found that it is:
    export GIN_MODE="release"
    export PORT=8080

    cd /srv/personal-site-go
    exec ./personal-site-go >> /var/log/personal-site-go.log 2>&1
end script
